# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:songlist]
command=gunicorn songbook.wsgi:application --workers 1 --bind=unix:/tmp/songlist.sock
directory=/home/dstarod/songbook
autorestart=true
redirect_stderr=true
stdout_logfile=/home/dstarod/songbook/songlist.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn